# common/hexagonrpc/hexagonrpcd-adsp-rootpd.service
[Unit]
Description=Daemon to support Qualcomm Hexagon ADSP virtual filesystem for RootPD
ConditionPathExists=/dev/fastrpc-adsp

[Service]
ExecStart=/usr/bin/hexagonrpcd -f /dev/fastrpc-adsp -d adsp
Restart=always
RestartSec=3
User=fastrpc
Group=fastrpc

[Install]
WantedBy=multi-user.target
